# rootfs/etc/services.d/dropbear/run
#!/usr/bin/with-contenv sh
# shellcheck shell=sh disable=SC2086

exec /usr/sbin/dropbear ${DROPBEAR_OPTIONS}
